# portal-test.service
[Unit]
Description=Portal test

[Service]
BusName=org.gnome.PortalTest
ExecStart=portal-test
Type=dbus
